# rootfs/var/run/user/1000/systemd/generator.late/xdg-desktop-autostart.target.wants/app-snap\x2duserd\x2dautostart@autostart.service
# Automatically generated by systemd-xdg-autostart-generator

[Unit]
Documentation=man:systemd-xdg-autostart-generator(8)
SourcePath=/etc/xdg/autostart/snap-userd-autostart.desktop
PartOf=graphical-session.target

Description=Snap user application autostart helper
After=graphical-session.target

[Service]
Type=exec
ExecStart=:/usr/bin/snap "userd" "--autostart"
Restart=no
TimeoutSec=5s
Slice=app.slice
